# templates/etc/rc.d/sshd-openbsd
#!/bin/sh
#
# OpenBSD rc.d script for <%= sshd_name %><%= suffix %>

daemon="<%= sshd_binary %><%= suffix %>"

. /etc/rc.d/rc.subr

rc_cmd $1
